# autossh/rootfs/etc/services.d/autossh/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting autossh server..."
exec /tmp/go.sh

exec nginx
